# cron/agatta.cron
#!/bin/sh


cd "/working/gdrive-local/Colab Notebooks/METIS-4.AI"


python agatta_tasks.py
